# immich/systemd/mnt-storage_pool.automount
[Unit]
Description=Automount MergerFS Storage Pool

[Automount]
Where=/mnt/storage_pool

[Install]
WantedBy=multi-user.target
